# paper2skills-vault/07-NLP-VOC/research/02-脚本工具/01-标签进化/com.momcozy.voc.monthly-evolution.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Phase 5 D12 T12.2 — macOS LaunchAgent for monthly VOC evolution

  安装：
    cp paper2skills-vault/07-NLP-VOC/research/02-脚本工具/01-标签进化/com.momcozy.voc.monthly-evolution.plist \
       ~/Library/LaunchAgents/

  启用：
    launchctl load ~/Library/LaunchAgents/com.momcozy.voc.monthly-evolution.plist

  立即触发一次（QA 场景 2）：
    launchctl kickstart -k gui/$(id -u)/com.momcozy.voc.monthly-evolution

  停用：
    launchctl unload ~/Library/LaunchAgents/com.momcozy.voc.monthly-evolution.plist

  日程：每月 1 日 02:00 触发（决策 3，§9）。
  飞书通知：脚本读取 ~/.paper2skills/feishu_webhook 自动推送完成/失败消息。
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.momcozy.voc.monthly-evolution</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/pray/project/paper_to_skills/paper2skills-vault/07-NLP-VOC/research/02-脚本工具/01-标签进化/monthly_evolution_cron.py</string>
        <string>--tenant</string>
        <string>momcozy</string>
        <string>--log-out</string>
        <string>/Users/pray/project/paper_to_skills/paper2skills-vault/07-NLP-VOC/research/04-输出结果/03-审计报告/monthly_evolution_latest.md</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/voc_monthly_evolution.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/voc_monthly_evolution.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/pray/project/paper_to_skills</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>
</dict>
</plist>
